# RemoteContent/Hg/0.0.4/hw.service
[Unit]
Description=hw
After=multi-user.target
[Service]
Type=simple
ExecStart=/bin/bash /var/gvg/rw.sh
 
[Install]
WantedBy=multi-user.target